# openwrt/package/yodaos/vol_ctrl/hisi_nau8540_ctrld
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org

START=89
STOP=2

start(){
  /usr/bin/hisitools/set_nau8540.sh
}

stop(){
	echo "stop"
}
