# openwrt/package/yodaos/vold/mount_data
#!/bin/sh /etc/rc.common

START=01

mount_data() {
  vold
}

start() {
  mount_data
  mkdir -p /data/upgrade
}
